# setup/email-indexer/email-indexer-incremental.service
[Unit]
Description=Email FTS5 indexer (incremental, mtime-based)
After=dovecot.service
Wants=dovecot.service

[Service]
Type=oneshot
User=root
Group=root
Nice=10
IOSchedulingClass=idle
ExecStart=/usr/bin/python3 /opt/email-indexer/indexer.py --incremental
StandardOutput=journal
StandardError=journal
TimeoutStartSec=900
# After each run, sync DB ownership so it's readable for rsync user
ExecStartPost=/bin/chmod 644 /var/lib/email-indexer/live.db
ExecStartPost=-/bin/chmod 644 /var/lib/email-indexer/live.db-wal
ExecStartPost=-/bin/chmod 644 /var/lib/email-indexer/live.db-shm
