# setup/email-indexer/email-indexer-incremental.service
[Unit]
Description=Email indexer (incremental, mtime-based) → PostgreSQL
After=dovecot.service postgresql.service
Wants=dovecot.service postgresql.service

[Service]
Type=oneshot
User=root
Group=root
Nice=10
IOSchedulingClass=idle
EnvironmentFile=/etc/mailinabox/postgres.env
Environment=PG_DSN=${PG_DSN_INDEXER}
ExecStart=/usr/bin/flock -w 300 /var/lock/email-indexer.lock /usr/bin/python3 /opt/email-indexer/indexer.py --incremental
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600
